# lib/systemd/system/networking.service.d/30_apparmor_profile_everything.conf
## Copyright (C) 2012 - 2021 ENCRYPTED SUPPORT LP <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

[Service]
ExecStart=
ExecStart=/usr/sbin/networking-aae up
ExecStop=
ExecStop=/usr/sbin/networking-aae down
