# lib/systemd/system/haveged.service.d/30_apparmor_profile_everything.conf
## Copyright (C) 2019 - 2022 ENCRYPTED SUPPORT LP <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

[Service]
## no_new_privs blocks transitions to haveged's AppArmor profile.
NoNewPrivileges=no
