# lib/systemd/system/haveged.service.d/30_apparmor_profile_everything.conf
[Service]
## no_new_privs blocks transitions to haveged's AppArmor profile.
NoNewPrivileges=no
